# infra/docker/builder/supervisord.conf
[supervisord]
# フォアグラウンドで起動
nodaemon=true
# ログ
logfile=/var/log/supervisord.log

[program:crond]
# 実行コマンド
command=/sbin/crond -n
# プログラムの自動起動
autostart=true
# プログラム停止後の再起動
autorestart=true